# airpi/package/feeds/luci/luci-app-Airpifanctrl/root/etc/init.d/Airpifanctrl
#!/bin/sh /etc/rc.common
START=99

printMsg() {
    local msg="$1"
    logger -t "${PROGRAM}" "${msg}"
} 

restart(){                                                                                                       
  /usr/bin/fancts.sh &
  printMsg "Airpi风扇控制模块开始工作！"
}

start(){
  /usr/bin/fancts.sh &
  printMsg "Airpi风扇控制模块开始工作！"
}
